# crontab
# ===== FULL US SCANS =====
# Every day at 1 AM - run full US scan without page limit to get all campgrounds
0 1 * * * curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"region": "us", "max_pages": 10}' >> /var/log/cron.log 2>&1

# Every Monday at 3 AM - run additional US bounds scan to ensure full coverage
0 3 * * 1 curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"region": "us", "max_pages": 10}' >> /var/log/cron.log 2>&1

# ===== PARALLEL SCAN - MULTIREGION =====
# Every day at 2 AM - run a parallel multiregion scan for faster and more thorough coverage
0 2 * * * curl -X POST http://api:8000/scrape-multiregion -H "Content-Type: application/json" -d '{"max_workers": 8, "max_pages": 10}' >> /var/log/cron.log 2>&1

# ===== ADDRESS UPDATE JOBS =====
# Update addresses for campgrounds without them - using parallel geocoding with 8 workers
0 */6 * * * curl -X POST http://api:8000/update-addresses -H "Content-Type: application/json" -d '{"limit": 200, "max_workers": 8}' >> /var/log/cron.log 2>&1

