# crontab
# ===== FULL US SCANS =====
# Every day at 1 AM - run full US scan (scan all regions systematically)
0 1 * * * curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"scrape_full_us": true, "max_pages": 5}' >> /var/log/cron.log 2>&1

# ===== MAJOR US REGIONS =====
# Every 12 hours - run US bounds scan
0 */12 * * * curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"region": "us", "max_pages": 8}' >> /var/log/cron.log 2>&1

# Every Monday, Wednesday, Friday at 2 AM - run Western US region
0 2 * * 1,3,5 curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"region": "western_us", "max_pages": 6}' >> /var/log/cron.log 2>&1

# Every Monday, Wednesday, Friday at 4 AM - run Eastern US region
0 4 * * 1,3,5 curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"region": "eastern_us", "max_pages": 6}' >> /var/log/cron.log 2>&1

# Every Tuesday, Thursday, Saturday at 2 AM - run Midwest US region
0 2 * * 2,4,6 curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"region": "midwest_us", "max_pages": 6}' >> /var/log/cron.log 2>&1

# Every Tuesday, Thursday, Saturday at 4 AM - run Southern US region
0 4 * * 2,4,6 curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"region": "southern_us", "max_pages": 6}' >> /var/log/cron.log 2>&1

# ===== POPULAR STATES =====
# Every Sunday at 3 AM - run California region
0 3 * * 0 curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"region": "california", "max_pages": 5}' >> /var/log/cron.log 2>&1

# Every Sunday at 5 AM - run New York region
0 5 * * 0 curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"region": "new_york", "max_pages": 5}' >> /var/log/cron.log 2>&1

# ===== TEST FOR GEOCODING =====
# Run in 5 minutes after container start to populate address fields
*/5 * * * * curl -X POST http://api:8000/scrape -H "Content-Type: application/json" -d '{"scrape_full_us": true, "max_pages": 1}' >> /var/log/cron.log 2>&1

